# configs/bdfs-scrub.timer
[Unit]
Description=Weekly BTRFS scrub (bdfs)
Documentation=https://gitlab.com/openos-project/linux-kernel_filesystem_deving/btrfs-dwarfs-framework

[Timer]
# Run every Sunday at 02:00 local time.
# RandomizedDelaySec spreads load across a fleet.
OnCalendar=Sun *-*-* 02:00:00
RandomizedDelaySec=30min
Persistent=true

[Install]
WantedBy=timers.target
